# root/etc/services.d/sabnzbd/run
#!/usr/bin/with-contenv bash

#variables set in /etc/cont-init.d/20 I think

exec \
	s6-setuidgid abc "$APP_COMP" "$APPDIR/${APP_EXEC}" \
	"$APP_OPTS"


